# cron/crontab
# Auto sync ip addresses of Cloudflare and reload nginx
# Every day at 8:30
30 8 * * * /usr/local/bin/update_cloudflare_ips.sh restart_nginx >/dev/null 2>&1
